# etc/cron.daily/archive_data
#!/bin/bash
python3 /home/gpstracker/GPSTracker/gps_tracker/archive_data.py
