# etc/systemd/system/barometer_poller.service
[Unit]
Description=barometer poller
After=redis-server.service

[Service]
Type=simple
User=pi
Restart=on-failure
ExecStart=/home/pi/GPSTracker/gps_tracker/barometer_poller.py

[Install]
WantedBy=multi-user.target
