# etc/systemd/system/gpspoller.service
[Unit]
Description=GPSD poller
After=redis-server.service

[Service]
Type=simple
User=gpstracker
Restart=on-failure
ExecStart=/home/gpstracker/GPSTracker/gps_tracker/gps_poller.py
WorkingDirectory=/home/gpstracker/GPSTracker/gps_tracker

[Install]
WantedBy=multi-user.target
