# etc/systemd/system/imu_baro_poller.service
[Unit]
Description=IMU and barometer poller
After=redis-server.service

[Service]
Type=simple
User=gpstracker
Restart=on-failure
ExecStart=/home/gpstracker/GPSTracker/gps_tracker/imu_baro_poller.py

[Install]
WantedBy=multi-user.target
